# postgres-dev/container-postgres-dev.service
# container-postgres-dev.service
# autogenerated by Podman 4.9.3
# Tue Jun  9 20:49:34 WIB 2026

[Unit]
Description=Podman container-postgres-dev.service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=%t/containers

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
TimeoutStopSec=70
ExecStart=/usr/bin/podman run \
	--cidfile=%t/%n.ctr-id \
	--cgroups=no-conmon \
	--rm \
	--sdnotify=conmon \
	--replace \
	--name=postgres-dev \
	-d \
	--label io.podman.compose.config-hash=46fcecccd95442029a36f47ce0228d7950706324af24df6910dae85070ff6533 \
	--label io.podman.compose.project=postgres-dev \
	--label io.podman.compose.version=1.0.6 \
	--label PODMAN_SYSTEMD_UNIT=podman-compose@postgres-dev.service \
	--label com.docker.compose.project=postgres-dev \
	--label com.docker.compose.project.working_dir=/home/imam/Documents/development/javascript-typescript/frayukti-dexiejs-app/postgres-dev \
	--label com.docker.compose.project.config_files=compose.yml \
	--label com.docker.compose.container-number=1 \
	--label com.docker.compose.service=postgres \
	-e POSTGRES_USER=appuser \
	-e POSTGRES_PASSWORD=apppassword \
	-e POSTGRES_DB=appdb \
	-v postgres-dev_pgdata:/var/lib/postgresql/data \
	--net postgres-dev_default \
	--network-alias postgres \
	-p 5432:5432 docker.io/library/postgres:16
ExecStop=/usr/bin/podman stop \
	--ignore -t 10 \
	--cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm \
	-f \
	--ignore -t 10 \
	--cidfile=%t/%n.ctr-id
Type=notify
NotifyAccess=all

[Install]
WantedBy=default.target
